# supervisord.conf
# https://github.com/Supervisor/supervisor/issues/376#issuecomment-404385767
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
# ---------

[supervisord]
logfile=/var/log/supervisor/supervisord.log
nodaemon=true
user=root
[program:apache]
command=httpd -D "FOREGROUND" -k start
[program:keycloak]
command=/opt/keycloak-4.8.3.Final/bin/standalone.sh -b 0.0.0.0
[program:startup.sh]
command=/opt/keycloak-4.8.3.Final/startup.sh -D